# roles/geoipupdate/templates/geoipupdate.service.j2
[Unit]
Description=Update MaxMind GeoIP databases

[Service]
Type=oneshot
ExecStart=/usr/bin/geoipupdate
TimeoutSec=600
