# roles/geoipupdate/templates/geoipupdate.timer.j2
[Unit]
Description=Update MaxMind GeoIP databases

[Timer]
OnCalendar=weekly
Persistent=yes

[Install]
WantedBy=timers.target
