# roles/zfs/templates/syncoid.timer.j2
[Unit]
Description=Sync ZFS datasets

[Timer]
OnCalendar={{ item.1.freq | default('daily') }}

[Install]
WantedBy=timers.target
