# roles/common/templates/rc.local.j2
#!/bin/bash

# {{ ansible_managed }}

{% for cmd in system_rc_local_cmd %}
{{ cmd }}
{% endfor %}

exit 0
