# roles/mailman/templates/mailman-digests.timer.j2
[Unit]
Description=Mailman digest sender

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
