# roles/unmaintained/psono/templates/psono-cleartoken.timer.j2
[Unit]
Description=Psono cleartoken

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
