# roles/zimbra/templates/zmldapsync.service.j2
[Unit]
Description=Sync LDAP accounts into Zimbra

[Service]
Type=oneshot
ExecStart=/opt/zimbra/bin/zmldapsync
TimeoutSec=300
