# roles/jitsi_jibri/templates/jitsi-jibri-cleaner.timer.j2
[Unit]
Description=Jibri Recordings Cleaner

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

