# roles/zabbix_lld_all_graph/templates/zabbix-lld-all-graph.timer.j2
[Unit]
Description=Update graphs with multiple lld items

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
