# roles/funkwhale/templates/funkwhale-update-media.timer.j2
[Unit]
Description=Update funkwhale media library

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
