# roles/wh_zimbra/templates/zmldapsync-wh.timer.j2
[Unit]
Description=Sync LDAP Users with Zimbra for hosted clients

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
