# roles/wh_zimbra/templates/zmldapsync-wh.service.j2
[Unit]
Description=Sync LDAP accounts into Zimbra for hosted clients

[Service]
Type=oneshot
ExecStart=/opt/zimbra/bin/zmldapsync --config /opt/zimbra/conf/zmldapsync-wh.yml
TimeoutSec=300
