# roles/seafile/templates/seafile-gc.timer.j2
[Unit]
Description=Seafile Garbage Collect service

[Timer]
OnCalendar=Sun *-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target
