# roles/seafile/templates/seafile-clean-db.service.j2
[Unit]
Description=Seafile DB cleaner service

[Service]
Type=oneshot
PrivateTmp=yes
ExecStart={{ seafile_root_dir }}/bin/clean_db.sh
