# roles/pmg/templates/imap-sa-learn.timer.j2
[Unit]
Description=Train spamassassin beayes from an IMAP account

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
